# scripts/meanmcq_super.conf
[program:meanmcq]
command=/home/red/IdeaProjects/meanMCQ/scripts/app.run
user=red
redirect_stderr=true
stdout_logfile=/home/red/IdeaProjects/meanMCQ/logs/meanmcq.log
stopasgroup=true
